# .docker/supervisor/supervisor.ini
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=php-fpm -F -O -d 'error_log = /dev/fd/2'
autorestart=unexpected
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off; error_log /dev/fd/2 info;'
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


